# supervisord.conf
[supervisord]
nodaemon=true

[program:django]
command=python manage.py runserver 0.0.0.0:$PORT
autostart=true
autorestart=true

[program:consumer]
command=python3 -m messaging.consumer
autostart=true
autorestart=true
